# resource/simple-service/simple.service
[Unit]
Description=Simple Test Service
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/simple-service
ExecStart=/usr/bin/python3 -m http.server 8080
Restart=always

[Install]
WantedBy=multi-user.target 